# systemd/restic-prune@.service
[Unit]
Description=Restic prune task %I
[Service]
Type=oneshot
ExecStart=restic prune
EnvironmentFile=/etc/restic/restic-backup.%i.conf
