# systemd/restic-backup@.timer
[Unit]
Description=Backup "%I" with restic daily
[Timer]
OnCalendar=*-*-* 1:00:00
Persistent=true
[Install]
WantedBy=timers.target
